# system/systemd/user/quickshell-network.service
[Unit]
Description=Quickshell Network Panel (WiFi + BT)
After=graphical-session.target
After=quickshell.service

[Service]
ExecStartPre=/usr/bin/sh -c 'echo -n > "$XDG_RUNTIME_DIR/qs-netpanel-cmd"'
ExecStart=/usr/bin/qs -p %h/.config/quickshell-network/shell.qml
StartLimitBurst=5
StartLimitIntervalSec=30
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
