# system/systemd/user/wave-check.service
[Unit]
Description=Quickshell wave video checker
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.config/quickshell/wave-check.sh
RemainAfterExit=yes

[Install]
WantedBy=default.target
